# docker/php/cron.d/cron
* * * * * cd /srv/api && php bin/console scheduler:consume
